# clack.service
[Unit]
Description=Clack - Mechanical Keyboard Click Sound App
After=graphical-session.target
Wants=graphical-session.target
PartOf=graphical-session.target
StartLimitIntervalSec=0

[Service]
Type=simple
ExecStart=/usr/bin/env python3 /home/konlan/Work/clack/clack.py
Restart=on-failure
RestartSec=2
Environment=QT_QPA_PLATFORM=wayland

[Install]
WantedBy=graphical-session.target
